# proxy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[program:iptv-dashboard]
command=gunicorn --bind 0.0.0.0:5000 --workers 4 --threads 2 dashboard.app:app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/iptv-dashboard.err.log
stdout_logfile=/var/log/iptv-dashboard.out.log

[program:stream-manager]
command=python3 /app/proxy/stream_manager.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/stream-manager.err.log
stdout_logfile=/var/log/stream-manager.out.log
